# apps/docker-platform/source/docker-platform-bootstrap.service
[Unit]
Description=Docker Platform first boot bootstrap
Wants=network-online.target docker.service
After=network-online.target docker.service cloud-init.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/docker-platform-bootstrap.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
