# apps/anythingllm/source/anythingllm-bootstrap.service
[Unit]
Description=AnythingLLM Docker Bootstrap
After=docker.service network-online.target
Wants=docker.service network-online.target
ConditionPathExists=/opt/anythingllm/docker-compose.yml

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/anythingllm-bootstrap.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
