# apps/woocommerce/source/woocommerce-bootstrap.service
[Unit]
Description=WooCommerce Docker bootstrap
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/woocommerce-bootstrap.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
